# sdvx_headless.service
[Unit]
Description=sdvx_controller

[Service]
Type=simple
ExecStart=/opt/sdvx_headless.py
[Install]
WantedBy=multi-user.target
